# root/etc/s6-overlay/s6-rc.d/laravel-scheduler/run
#!/bin/sh

# Change to the application directory
cd /var/www/html

# Run the Laravel scheduler every minute
php artisan schedule:work
